# cron/crontab
#crontab
0-59/1  *  *  *  * root *applications/general_ledger/modules/crontest.py
