# cron/crontab
#crontab
0-59/1  *  *  *  * root *applications/instore/modules/background_updater.py
